# tools/heat_test.conf
description "RPi Heat Test"
author "born.smasher"

start on runlevel [2345]
stop on runlevel [!2345]

chdir /opt/RpiTemp
exec python heat_test.py

respawn
